# univention-ucc-pam/ucc-nss.conf
# ucc - nss update 

description     "updating ucc nss information"

start on net-device-up

task

script
	univention-ucc-update-nss
end script
