# univention-corporate-client/ucc-update.conf
# UCC - Install software updates and additional packages

description     "Installing software updates"

start on net-device-up IFACE=eth0 and stopped ucc-syspol

task

exec	/usr/sbin/univention-ucc-software-update

