# univention-corporate-client/conffiles/etc/profile.d/ucc.sh
@%@UCRWARNING=# @%@

# Expand the $PATH to include /sbin so initctl can be started
# (required by ubuntu desktop session)
PATH=$PATH:/sbin:/usr/sbin
