# univention-corporate-client/conffiles/etc/update-motd.d/10-help-text
#!/bin/sh
# The original script displays ubuntu specific help infos
