# univention-corporate-client/conffiles/etc/update-motd.d/90-updates-available
#!/bin/sh
# The original script runs a complete apt-get update, which is handled through univention-ucc-software-update in UCC
