# univention-corporate-client/conffiles/etc/update-motd.d/91-release-upgrade
#!/bin/sh
# The original script installs new Ubuntu updates, not used in UCC
